# ansible/roles/profiler/files/pprotein-agent.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/pprotein-agent

[Install]
WantedBy=multi-user.target
